# deploy/systemd/bestcard-ingest.service.d/override.conf
# Drop-in hardening for bestcard-ingest.service (CT 104).
# Install on the server as:
#   /etc/systemd/system/bestcard-ingest.service.d/override.conf
# then: systemctl daemon-reload && systemctl restart bestcard-ingest
#
# Fixes the OOM described in docs/incidents/2026-07-20-ingestion-oom.md.
[Service]
# Root cause fix: ThreadingHTTPServer spawns a thread per request and glibc
# gives each thread its own malloc arena (default up to 8*CPU), each hoarding
# freed large buffers (PDF/HTML/LLM text). Capping arenas to 2 keeps RSS flat.
Environment=MALLOC_ARENA_MAX=2
# Return freed heap to the OS sooner (bytes). Pairs with malloc_trim() in code.
Environment=MALLOC_TRIM_THRESHOLD_=131072

# Belt-and-suspenders ceiling on THIS service only: throttle at 600M, hard cap
# at 768M. If it ever exceeds the cap, systemd restarts just the service — the
# container (and SearXNG/Valkey) stay healthy instead of the whole box OOMing.
MemoryHigh=600M
MemoryMax=768M

# Always recover automatically.
Restart=always
RestartSec=2
